# infra/systemd/nexostrat-mirror-codeberg.path
[Unit]
Description=Nexostrat — watch Gitea bare-repo refs/heads/main and trigger Codeberg mirror
Documentation=file:///srv/Nexostrat/00_GOVERNANCE/system_map.md

[Path]
PathChanged=/srv/gitea/data/git/repositories/nexostrat/nexostrat.git/refs/heads/main
Unit=nexostrat-mirror-codeberg.service

[Install]
WantedBy=multi-user.target
