# infra/systemd/nexostrat-mirror-github.service
[Unit]
Description=Nexostrat — push HP origin to GitHub mirror
Documentation=file:///srv/Nexostrat/00_GOVERNANCE/system_map.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ricardo
Group=ricardo
WorkingDirectory=/srv/Nexostrat
Environment=HOME=/home/ricardo
Environment=SSH_AUTH_SOCK=
ExecStart=/srv/Nexostrat/infra/scripts/mirror-push.sh github

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/srv/Nexostrat
ProtectHome=read-only
# ricardo's HOME needs read for ~/.ssh/config + ~/.ssh/nexostrat_ed25519
BindReadOnlyPaths=/home/ricardo/.ssh

# Modest log retention
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexostrat-mirror-github

[Install]
WantedBy=multi-user.target
